# rustic/.config/systemd/user/rustic-backup.timer
[Unit]
Description=Her 4 saatte bir Rustic yedekleme

[Timer]
OnUnitActiveSec=4h
Persistent=true

[Install]
WantedBy=default.target
