# rustic/.config/systemd/user/rustic-backup.service
[Unit]
Description=Rustic user-level backup service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rustic-backup.sh

