# rustic/.config/systemd/user/rustic-backup.timer
[Unit]
Description=Her 4 saatte bir Rustic yedekleme

[Timer]
OnBootSec=5min
OnUnitActiveSec=8h
Persistent=true

[Install]
WantedBy=default.target
